# deploy/raspy-monitor.service
[Unit]
Description=Raspy Monitor — Raspberry Pi system monitoring dashboard
Documentation=https://github.com/K1taru/pi-monitor

[Service]
Type=simple
User=k1taru
Group=k1taru
WorkingDirectory=/home/k1taru/pi-monitor/backend
EnvironmentFile=/home/k1taru/pi-monitor/backend/.env
ExecStart=/home/k1taru/pi-monitor/backend/venv/bin/python app.py

# Restart on failure with a short delay
Restart=on-failure
RestartSec=5s

# Logging goes to journald — view with: journalctl -u raspy-monitor -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raspy-monitor

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/k1taru/pi-monitor/backend

[Install]
WantedBy=multi-user.target
