# deploy/raspy-monitor.service
[Unit]
Description=Raspberry Pi Monitor Backend
After=network.target

[Service]
Type=simple
User=k1taru
WorkingDirectory=/home/k1taru/pi-monitor/backend

# Load SECRET_KEY, JWT_SECRET_KEY, PORT, CORS_ORIGINS
EnvironmentFile=/home/k1taru/pi-monitor/backend/.env

# Use the venv interpreter so all pip packages are available
ExecStart=/home/k1taru/pi-monitor/backend/venv/bin/python app.py

Restart=always
RestartSec=5

# NOTE: Do NOT add NoNewPrivileges=yes here.
# The backend delegates privileged writes (fan PWM, CPU governor, reboot) to
# root-owned wrapper scripts via passwordless sudo (see /etc/sudoers.d/raspy-monitor).
# NoNewPrivileges would block sudo from gaining root and break those calls.

[Install]
WantedBy=multi-user.target
