# setup/pi-monitor.service
# pi-monitor.service — Reference template
#
# This file shows the structure of the service unit that setup.sh generates.
# setup.sh writes the actual file to /etc/systemd/system/pi-monitor-<username>.service
# with all variables filled in from backend/.env and the username argument.
#
# You do NOT need to edit or install this file manually.
# Run:  sudo setup/setup.sh <username>

[Unit]
Description=Pi Monitor Backend
After=network.target

[Service]
Type=simple
User=<username>
WorkingDirectory=<PROJECT_DIR>/backend
EnvironmentFile=<PROJECT_DIR>/backend/.env
ExecStart=<PROJECT_DIR>/backend/venv/bin/python app.py
Restart=always
RestartSec=5

# NoNewPrivileges must NOT be set — the service uses passwordless sudo
# to call pi-monitor-fan-control and pi-monitor-gov-control.

[Install]
WantedBy=multi-user.target
